# config/crontab
# =============================================
#  Misskey Backup Scheduler
#  Runs database backup at specified intervals
# =============================================
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Run at 3:00 and 15:00 every day
0 3,15 * * * root /usr/local/bin/misskey-backup >> /var/log/cron.log 2>&1

# 空行が必要